# xray-tproxy-network.service
[Unit]
Description=Xray Transparent Proxy Network Configuration
Documentation=https://github.com/XTLS/Xray-core
After=network-online.target
Wants=network-online.target
Before=xray.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Setup IPv4 routing rules
ExecStartPre=/bin/bash -c 'ip rule show | grep -q "fwmark 0x1 lookup 100" || ip rule add fwmark 1 table 100'
ExecStartPre=/bin/bash -c 'ip route show table 100 | grep -q "local default" || ip route add local 0.0.0.0/0 dev lo table 100'

# Setup IPv6 routing rules
ExecStartPre=/bin/bash -c 'ip -6 rule show | grep -q "fwmark 0x1 lookup 100" || ip -6 rule add fwmark 1 table 100'
ExecStartPre=/bin/bash -c 'ip -6 route show table 100 | grep -q "local default" || ip -6 route add local ::/0 dev lo table 100'

# Load nftables rules
ExecStart=/usr/sbin/nft -f /etc/nftables/xray-tproxy.nft

# Reload: reapply all rules
ExecReload=/bin/bash -c 'ip rule show | grep -q "fwmark 0x1 lookup 100" || ip rule add fwmark 1 table 100'
ExecReload=/bin/bash -c 'ip route show table 100 | grep -q "local default" || ip route add local 0.0.0.0/0 dev lo table 100'
ExecReload=/bin/bash -c 'ip -6 rule show | grep -q "fwmark 0x1 lookup 100" || ip -6 rule add fwmark 1 table 100'
ExecReload=/bin/bash -c 'ip -6 route show table 100 | grep -q "local default" || ip -6 route add local ::/0 dev lo table 100'
ExecReload=/usr/sbin/nft -f /etc/nftables/xray-tproxy.nft

# Cleanup on stop
ExecStop=/bin/bash -c 'ip rule del fwmark 1 table 100 2>/dev/null || true'
ExecStop=/bin/bash -c 'ip route del local 0.0.0.0/0 dev lo table 100 2>/dev/null || true'
ExecStop=/bin/bash -c 'ip -6 rule del fwmark 1 table 100 2>/dev/null || true'
ExecStop=/bin/bash -c 'ip -6 route del local ::/0 dev lo table 100 2>/dev/null || true'
ExecStop=/bin/bash -c 'nft delete table inet xray 2>/dev/null || true'

# Security settings
NoNewPrivileges=true
ProtectSystem=false
ProtectHome=true
PrivateTmp=true

# Restart settings
TimeoutStartSec=30
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
